# profiles/loki/uci-defaults/90-network-ip
# Delete the original device record (br-lan bridge)
uci delete network.@device[0]

# Delete the original lan interface (lan interface)
uci delete network.lan

# Configure WAN interface
uci set network.wan=interface
uci set network.wan.proto=dhcp
uci set network.wan.device='sw-eth1'

# Create a new device and configure it as a bridge
uci set network.fry_bridge=device
uci set network.fry_bridge.name='fry-bridge'
uci set network.fry_bridge.type='bridge'
uci add_list network.fry_bridge.ports='sw-eth2'

# Configure VLANs on the bridge
uci set network.vlan1=bridge-vlan
uci set network.vlan1.device='fry-bridge'
uci set network.vlan1.vlan='1'
uci add_list network.vlan1.ports='sw-eth2:t'

uci set network.vlan20=bridge-vlan
uci set network.vlan20.device='fry-bridge'
uci set network.vlan20.vlan='20'
uci add_list network.vlan20.ports='sw-eth2'

# Configure the captive interface
uci set network.captive=interface
uci set network.captive.proto='static'
uci set network.captive.ipaddr='192.168.55.1'
uci set network.captive.netmask='255.255.255.0'
uci set network.captive.device='fry-bridge.1'
uci set network.captive.ip6assign='60'

# Configure the encrypted interface
uci set network.encrypted=interface
uci set network.encrypted.proto='static'
uci set network.encrypted.ipaddr='192.168.45.1'
uci set network.encrypted.netmask='255.255.255.0'
uci set network.encrypted.device='fry-bridge.20'
uci set network.encrypted.ip6assign='60'

del network.lan

# Commit changes and delete the script to avoid re-execution
uci commit network
