# profiles/artemisa/uci-defaults/90-dhcp-config
uci -q batch << EOI
del dhcp.lan.ra_slaac
commit dhcp
EOI