# profiles/odyssey/uci-defaults/99-radsec-config
rm /etc/config/radsecproxy
touch /etc/config/radsecproxy

uci -q batch << EOI

add radsecproxy options
rename radsecproxy.@options[0]=general

add_list radsecproxy.general.ListenUDP="127.0.0.1:1812"
add_list radsecproxy.general.ListenUDP="127.0.0.1:1813"
set radsecproxy.general.LogLevel="1"
set radsecproxy.general.LogDestination="x-syslog:///"

add radsecproxy tls
rename radsecproxy.@tls[0]=tls

set radsecproxy.tls.CACertificatePath="/etc/fry-os-services/data"
set radsecproxy.tls.CACertificateFile="/etc/fry-os-services/data/radsec-ca.crt"
set radsecproxy.tls.certificateFile="/etc/fry-os-services/data/radsec.crt"
set radsecproxy.tls.certificateKeyFile="/etc/fry-os-services/data/radsec.key"
set radsecproxy.tls.certificateKeyPassword=""
set radsecproxy.tls.name="tls_fry_openroaming"

add radsecproxy client
rename radsecproxy.@client[0]=client

set radsecproxy.client.name="localhost"
set radsecproxy.client.type="udp"
set radsecproxy.client.secret="local-proxy-secret"

add radsecproxy server
rename radsecproxy.@server[0]=fry_openroaming_upstream

set radsecproxy.fry_openroaming_upstream.name="128.203.161.64"
set radsecproxy.fry_openroaming_upstream.type="tls"
set radsecproxy.fry_openroaming_upstream.secret="radsec"
set radsecproxy.fry_openroaming_upstream.port="2083"
set radsecproxy.fry_openroaming_upstream.tls="tls_fry_openroaming"
set radsecproxy.fry_openroaming_upstream.statusServer="1"
set radsecproxy.fry_openroaming_upstream.certificateNameCheck="0"
set radsecproxy.fry_openroaming_upstream.loopPrevention="1"
set radsecproxy.fry_openroaming_upstream.tcpKeepalive="0"

add radsecproxy realm
rename radsecproxy.@realm[0]=fry_openroaming_realms

set radsecproxy.fry_openroaming_realms.name="*"
add_list radsecproxy.fry_openroaming_realms.server="128.203.161.64"
add_list radsecproxy.fry_openroaming_realms.accountingServer="128.203.161.64"

commit radsecproxy


EOI
