# profiles/odyssey/uci-defaults/99-wireless-operator
uci -q batch << EOI

set wireless.radio0.disabled='0'
set wireless.radio0.htmode='HT40'
set wireless.radio0.channel='auto'
set wireless.radio0.cell_density='0'

set wireless.default_radio0.disabled='0'
set wireless.default_radio0.ssid='Fry Operator'
set wireless.default_radio0.encryption='psk2'
set wireless.default_radio0.key='988Dh283zX'
set wireless.default_radio0.network='lan'
set wireless.default_radio0.isolate='1'

rename wireless.default_radio0=default_wifi_interface_0

set wireless.wifinet2=wifi-iface

set wireless.wifinet2.disabled='0'
set wireless.wifinet2.device='radio0'
set wireless.wifinet2.encryption='wpa2'
set wireless.wifinet2.key='local-proxy-secret'
set wireless.wifinet2.network='lan'
set wireless.wifinet2.mode='ap'
set wireless.wifinet2.port='1812'
set wireless.wifinet2.ssid='Fry Mobile'
set wireless.wifinet2.isolate='1'
set wireless.wifinet2.auth_server='127.0.0.1'
set wireless.wifinet2.auth_port='1812'
set wireless.wifinet2.auth_secret='local-proxy-secret'
set wireless.wifinet2.server='127.0.0.1'
set wireless.wifinet2.acct_server='127.0.0.1'
set wireless.wifinet2.acct_port='1813'
set wireless.wifinet2.acct_secret='local-proxy-secret'
set wireless.wifinet2.acct_interval='300'
set wireless.wifinet2.dae_port='3799'
set wireless.wifinet2.nasid='40:A5:EF:B9:B8:3F'
set wireless.wifinet2.hs20='1'
set wireless.wifinet2.hs20_oper_friendly_name='eng:Fry'
set wireless.wifinet2.hs20_operating_class='5173'
set wireless.wifinet2.ieee80211w='0'
set wireless.wifinet2.ifname='wlan0-2'
set wireless.wifinet2.iw_enabled='1'
set wireless.wifinet2.iw_interworking='1'
set wireless.wifinet2.iw_internet='1'
set wireless.wifinet2.iw_venue_group='1'
set wireless.wifinet2.iw_venue_type='7'
set wireless.wifinet2.iw_uesa='1'
set wireless.wifinet2.iw_asra='0'
set wireless.wifinet2.iw_esr='1'
set wireless.wifinet2.iw_access_network_type='2'
set wireless.wifinet2.iw_disable_dgaf='1'
set wireless.wifinet2.iw_ipaddr_type_availability='11'
set wireless.wifinet2.iw_network_auth_type='00'
set wireless.wifinet2.iw_hessid='00:00:00:01:02:03'
add_list wireless.wifinet2.iw_domain_name='fry.network'
add_list wireless.wifinet2.iw_domain_name='openroaming.fry.network'
add_list wireless.wifinet2.iw_domain_name='openroaming.org'
add_list wireless.wifinet2.iw_domain_name='apple.openroaming.net'
add_list wireless.wifinet2.iw_domain_name='google.openroaming.net'
add_list wireless.wifinet2.iw_domain_name='ciscooneid.openroaming.net'
add_list wireless.wifinet2.iw_domain_name='samsung.openroaming.net'
add_list wireless.wifinet2.iw_nai_realm='0,fry.network,13[5:6],21[2:4][5:7]'
add_list wireless.wifinet2.iw_venue_name='eng:Fry'
add_list wireless.wifinet2.iw_venue_url='1:https://account.fry.network'

rename wireless.wifinet2=passpoint_fry_wifi_2ghz

commit wireless

EOI
