# profiles/voyage-installer/uci-defaults/99-wireless-operator
uci -q batch << EOI

uci set wireless.radio0=wifi-device
uci set wireless.radio0.type='mac80211'

uci set wireless.default_radio0=wifi-iface
uci set wireless.default_radio0.device='radio0'
uci set wireless.default_radio0.network='lan'
uci set wireless.default_radio0.mode='ap'  
uci set wireless.default_radio0.ssid='Fry Operator'
uci set wireless.default_radio0.encryption='psk2'
uci set wireless.default_radio0.key='988Dh283zX'
uci commit wireless


EOI