# profiles/chronos/uci-defaults/90-network-ip
uci -q batch << EOI
set network.wan=interface
set network.wan.proto='dhcp'
set network.wan.device='lan'
del network.lan.ip6assign
set network.lan.proto='static'
del network.lan.device
set network.lan.ipaddr='192.168.45.1'
set network.lan.netmask='255.255.255.0'
set network.lan.type='bridge'
commit network
EOI