# profiles/prometheus-installer/uci-defaults/90-dhcp-config
uci -q batch << EOI
uci del dhcp.lan.ra_slaac
commit dhcp
EOI