# profiles/genesis/uci-defaults/90-dropbear-keys
rm /etc/dropbear/authorized_keys
touch /etc/dropbear/authorized_keys
printf 'ssh-rsa AAAAB3NzaC1yc2EAAAADAQABAAABAQC4Y2vORpbLfqfwdKeyntKLFSiEEVUzXHCL4tBEMEnjMat/VYvvZfoh9+ZmcZs68h2FisRNXkMDHfUawEvm74z0a6sBnOUZ5prNHbwAF5QucFP4iEo98gMF7jxNTuZL0UVloDHJUK1frWEuZKQMXmXQCT9qYnHr/N5Lpxgm/NliPf1069LeQRCfpyRu7fwJHIxUdDdthNMGZ8lQsTWgdtIuvCUhevzZUqpiaDSZd00QPwQqzy4xdD7ZhOmSD6o9c0lm3gG3Mv9Pm3x1LshWi1aGfj0zlMkwcV/0WQJ9/Du8A24uKD6xOSJOIjXxATWVX+NnbwSMUCmHUOsoJt6xdtZZ root' >> /etc/dropbear/authorized_keys
chmod 600 /etc/dropbear/authorized_keys