# profiles/genesis/uci-defaults/99-wireless-operator
uci -q batch << EOI

set wireless.radio0.disabled='0'
set wireless.radio0.htmode='VHT20'
set wireless.radio0.channel='auto'
set wireless.radio0.txpower='23'
set wireless.radio0.cell_density='0'
set wireless.radio0.band='5g'

set wireless.radio1.disabled='0'
set wireless.radio1.htmode='HT20'
set wireless.radio1.channel='auto'
set wireless.radio1.txpower='27'
set wireless.radio1.cell_density='0'
set wireless.radio1.band='2g'

set wireless.default_radio0.disabled='0'
set wireless.default_radio0.ssid='Wayru Operator'
set wireless.default_radio0.encryption='sae-mixed'
set wireless.default_radio0.key='988Dh283zX'
set wireless.default_radio0.network='encrypted'
set wireless.default_radio0.isolate='1'

set wireless.default_radio1.disabled='0'
set wireless.default_radio1.ssid='Wayru Operator'
set wireless.default_radio1.encryption='sae-mixed'
set wireless.default_radio1.key='988Dh283zX'
set wireless.default_radio1.network='encrypted'
set wireless.default_radio1.isolate='1'

rename wireless.default_radio0=default_wifi_interface_0
rename wireless.default_radio1=default_wifi_interface_1

commit wireless

EOI