# profiles/apocalypse-installer/uci-defaults/90-network-ip
uci -q batch << EOI
set network.lan.ipaddr='192.168.45.1'
commit network
EOI